# boot2kube/rootfs-overlay/etc/systemd/system/kube-init.service
[Unit]
Description=kube init script

[Service]
ExecStart=/usr/sbin/kube-init.sh
Type=oneshot
RemainAfterExit=true
